# deploy/systemd/polar-feeder.service
[Unit]
Description=Polar Feeder BLE Controller
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=arcticproject
WorkingDirectory=/home/arcticproject/Desktop/polar_bear_project_github/polar-feeder
Environment=PYTHONPATH=/home/arcticproject/Desktop/polar_bear_project_github/polar-feeder/src/pi
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

ExecStartPre=/usr/bin/pinctrl set 17 op dl
ExecStartPre=/usr/bin/pinctrl set 27 op dl
ExecStart=/bin/bash -lc 'cd /home/arcticproject/Desktop/polar_bear_project_github/polar-feeder && exec bash src/pi/scripts/run.sh --ble-test --config /etc/polar_feeder/config.json'

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
